# templates/claude-gateway.service
[Unit]
Description=Claude AI Telegram Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=REPLACE_USER
Group=REPLACE_USER
WorkingDirectory=/home/REPLACE_USER/claude-gateway
ExecStart=/home/REPLACE_USER/claude-gateway/.venv/bin/python /home/REPLACE_USER/claude-gateway/gateway.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-gateway
Environment=HOME=/home/REPLACE_USER
Environment=PATH=/home/REPLACE_USER/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
